# home/Library/LaunchAgents/dev.dotfiles.update.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.dotfiles.update</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{ .chezmoi.homeDir }}/.local/bin/dotfiles-update</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>{{ .chezmoi.homeDir }}/.local/bin:{{ .homebrewPrefix }}/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>TZ</key>
        <string>Asia/Seoul</string>
    </dict>
    <!-- Fixed monthly schedule: day 1 and 16 at 14:00 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Day</key><integer>1</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Day</key><integer>16</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>0</integer></dict>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
